# share-lan_1.0_all/DEBIAN/postinst
#!/bin/sh
set -e

# Reload systemd so it picks up the new unit
if [ -d /run/systemd/system ]; then
    systemctl daemon-reload >/dev/null 2>&1 || true
    systemctl enable share-lan.service >/dev/null 2>&1 || true
    systemctl start share-lan.service >/dev/null 2>&1 || true
fi

exit 0
