# share-lan_1.0_all/DEBIAN/prerm
#!/bin/sh
set -e

if [ -d /run/systemd/system ]; then
    systemctl stop share-lan.service >/dev/null 2>&1 || true
    systemctl disable share-lan.service >/dev/null 2>&1 || true
fi

exit 0
